# Тестовое задание для Effective Mobile/testeffectmobile2.service
[Service]
ExecStart=/home/rinat/Interview/testeffectmobile2.sh

[Install]
WantedBy=multi-user.target